# ansible/playbooks/templates/minecraft.service.j2
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=simple
User={{ minecraft_user }}
Group={{ minecraft_user }}
WorkingDirectory={{ minecraft_dir }}

# Use the startup script instead of direct java command
ExecStart=/opt/minecraft/scripts/start-minecraft.sh

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=120

# Restart policy
Restart=on-failure
RestartSec=30

# Security settings
NoNewPrivileges=true

# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target